# rootfs/etc/services.d/dataannotation_projects/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting DataAnnotation Projects app..."
exec node -r tsx/cjs /app/src/main.ts
